# examples/release-post-2026-03-13.cron
CRON_TZ=America/New_York
18 9 13 3 * /bin/bash -lc 'set -euo pipefail; TODAY="$(date +\%F)"; [ "$TODAY" = "2026-03-13" ] || exit 0; LOCK_FILE="${X_RELEASE_POST_LOCK_FILE:-$HOME/.config/x-api/release-post-2026-03-13.lock}"; [ ! -e "$LOCK_FILE" ] || exit 0; OGK_REPO_DIR="${OGK_REPO_DIR:-$HOME/obsidian-growth-kit}"; MEMEGEN_PATH="${X_MEMEGEN_PATH:-$HOME/.openclaw/workspace-obsidian/out/memegen.py}"; MEME_PYTHON="${X_MEMEGEN_PYTHON:-$HOME/.openclaw/workspace-obsidian/out/.venv/bin/python}"; MEME_PATH="$(xmeme generate --generator "$MEMEGEN_PATH" --python "$MEME_PYTHON" --template charlie_fire.jpg --top-text "WHEN THE RELEASE POST IS FINALLY LIVE" --bottom-text "AND IT IS A FRESH MEME THIS TIME" --output-dir "$OGK_REPO_DIR/out" --prefix release-post)"; xpost --text-file "$OGK_REPO_DIR/x-post-release-package-2026-03-12.md" --media "$MEME_PATH"; mkdir -p "$(dirname "$LOCK_FILE")"; : > "$LOCK_FILE"'
